# .configs/supervisor/bagisto.conf
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx-error.log

# [program:varnish]
# command=/usr/sbin/varnishd -F -a :8081 -T :6082 -f /etc/varnish/default.vcl -s malloc,256m -S /etc/varnish/secret -n /var/lib/varnish 
# autostart=true
# autorestart=true
# stdout_logfile=/var/log/supervisor/varnish.log
# stderr_logfile=/var/log/supervisor/varnish-error.log
